# docker-assets/etc/service/log-forwarder/run
#!/bin/sh
exec tail -F -n 0 /var/log/gitlab-remote-mirror.log